# systemd/tests/ssh.service
[Unit]
Description=OpenBSD Secure Shell server
[Service]
ExecStart=/usr/sbin/sshd -D
[Install]
WantedBy=multi-user.target
